# conf/supervisor/bootsrv.ini
[program:bootsrv]
command=/bin/bash -c "\cp -rf  /srv/* /usr/local/webserver/coffee/web"
autostart=true
autorestart=false