# conf/supervisor/java.ini
[program:java]
command=/bin/bash -c "cd /usr/local/webserver/coffee/web && /usr/local/webserver/java/jre1.7.0_80/bin/java -jar website-1.0-SNAPSHOT.jar"
