# conf/supervisor/mysql.ini
[program:mysql]
command=/bin/bash -c "/usr/local/webserver/mysql/bin/mysqld_safe --user=mysql"
