# conf/supervisor/nginx.ini
[program:nginx]
command=/bin/bash -c "/usr/local/webserver/nginx/sbin/nginx"
