# conf/supervisor/nginx.ini
[program:nginx]
command=/bin/bash -c "/usr/local/webserver/nginx/sbin/nginx"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/worker-nginx.log