# conf/supervisor/redis.ini
[program:redis]
command=/bin/bash -c "/etc/init.d/redis start"
