# conf/supervisor/redis.ini
[program:redis]
command=/bin/bash -c "/etc/init.d/redis start"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/worker-redis.log